# deploy/io.rampant.photo-keywords.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template. install.sh substitutes __PYTHON__, __PROJECT__, __HOME__.
     Installed to ~/Library/LaunchAgents/ (per-user agent; Photos needs a GUI session). -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.rampant.photo-keywords</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__PROJECT__/src/run.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>2</integer>
        <key>Minute</key><integer>15</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>__PROJECT__</string>

    <!-- launchd gives an empty environment; set HOME so claude (~/.claude) and the
         secrets/config resolve, and a PATH covering ~/.local/bin (claude) plus the
         system dirs (/usr/bin/sips). The OAuth token is read by run.py from the
         secrets file, NOT placed here. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key><string>__HOME__</string>
        <key>PATH</key><string>__HOME__/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/state/photo-keywords/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.local/state/photo-keywords/logs/launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
